# images/redis/s6/services.d/exporter/run
#!/usr/bin/with-contenv sh

./redis_exporter
status=$?
if [ $status -ne 0 ]; 
then
  echo "Failed to start exporter: $status"
  exit $status
else
  echo "redis exporter started"
fi
