# images/redis/s6/services.d/exporter/run
#!/usr/bin/with-contenv sh
echo "starting redis exporter"
exec redis_exporter

