# images/redis/s6/services.d/redis/run
#!/usr/bin/with-contenv sh

redis-server
status=$?
if [ $status -ne 0 ]; then
  echo "Failed to start redis: $status"
  exit $status
else
  echo "redis started"
fi
